# files/installer/units/installer.service
[Unit]
Description=Bluefin Server DDI Installer Service
Documentation=https://github.com/projectbluefin/server
DefaultDependencies=no
ConditionPathExists=/opt/knuckle
After=basic.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/knuckle --os bluefin-ddi
StandardInput=tty
StandardOutput=tty
StandardError=journal+console
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes

[Install]
WantedBy=installer.target
