# files/installer/units/installer.service
[Unit]
Description=Bluefin Server DDI Installer Service
Documentation=https://github.com/projectbluefin/server
DefaultDependencies=no
ConditionPathExists=/opt/knuckle
After=basic.target systemd-udev-settle.service
Wants=systemd-udev-settle.service

[Service]
Type=idle
RemainAfterExit=no
ExecStartPre=modprobe nvme
ExecStartPre=modprobe nvme_core
ExecStartPre=udevadm settle --timeout=15
ExecStartPre=/bin/sh -c 'for i in $(seq 1 15); do [ -e /dev/disk/by-partlabel/bluefin-installer-data ] && exit 0; sleep 1; done; echo "DDI partition not found"; exit 1'
ExecStart=/usr/sbin/bluefin-autoinstall
StandardInput=tty
StandardOutput=tty
StandardError=journal+console
TTYPath=/dev/console
TTYReset=yes
TTYVHangup=yes

[Install]
WantedBy=installer.target
