# files/installer/units/installer.target
[Unit]
Description=Bluefin Server Installer
Documentation=https://github.com/projectbluefin/server
DefaultDependencies=no
# installer.service does the work; this target waits for it.
Wants=installer.service
After=installer.service
AllowIsolate=yes
