# files/installer/units/installer.service
[Unit]
Description=Bluefin Server DDI Installer Service
Documentation=https://github.com/projectbluefin/server
DefaultDependencies=no
ConditionPathExists=/usr/bin/bluefin-install
SuccessAction=poweroff

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/bluefin-install
FailureAction=none

[Install]
WantedBy=installer.target
