# files/os/systemd/systemd-sysupdate.service.d/kured-hook.conf
[Service]
ExecStartPost=/usr/bin/sh -c 'touch /run/reboot-required'
